# crontab
1,11,21,31,41,51  *  *  *  * cd /var/www/html/ && php -f downloader.php >/dev/null 2>&1
